# mn.emailer.conf
description "MyNumber emailer service"
author      "easanka@MyNumber"

# Used to Be: Start on Startup
# until we found some mounts weren't ready yet while booting:

start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/root"

    [ $(exec /usr/bin/php -f /var/html/mn.emailer.service/index.php) = 'ERROR' ] && ( stop; exit 1; )

end script

post-start script
   # Optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "mn.emailer.service has started!"
end script